# mysql_ansible/roles/make_mha_go/templates/mha-manager.service.j2
# Generated by Ansible (dbbot make_mha_go role). Do not edit manually.
[Unit]
Description=dbbot MHA Go Manager — {{ mha_go_cluster_name }}
Documentation=https://github.com/your-org/dbbot_mha
After=network.target {{ mysql_service_name }}.service
Wants={{ mysql_service_name }}.service

[Service]
Type=simple
User={{ mysql_user }}
Group={{ mysql_group }}
ExecStart={{ mha_go_binary_dest }} manager \
    --config {{ mha_go_config_dir }}/cluster.yaml \
    --log-format json
Restart=on-failure
RestartSec=10s
# Prevent fast restart loops
StartLimitInterval=60s
StartLimitBurst=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mha-manager

[Install]
WantedBy=multi-user.target
